# installer/system_files/shared/etc/skel/.config/autostart/on_gui_login.desktop
[Desktop Entry]
Type=Application
Name=Bienvenue dans Gablue Live
Exec=/usr/bin/on_gui_login.sh
OnlyShowIn=KDE;
X-KDE-autostart-phase=1
NoDisplay=true
