# files/system/all/usr/lib/systemd/system/plasmalogin.service.d/90-gablue-settle.conf
[Service]
# Remplace le settle udev aveugle de 10s par une attente ciblee GPU + input
# TEMPORAIRE : a supprimer quand Fedora/KDE corrige plasmalogin
ExecStartPre=
ExecStartPre=-/usr/libexec/gablue-wait-devices
