# config/files/usr/etc/profile.d/customperso.sh
# Personnal alias

alias goarch="distrobox enter archgab"

# Ignore duplicate
HISTCONTROL=ignoredups

