# config/files/usr/etc/profile.d/customperso.sh
# Personnal alias

alias goarch="distrobox enter archgab"

alias goarch-root="distrobox enter --root archgab"

# Ignore duplicate
HISTCONTROL=ignoredups

