# files/system/all/etc/xdg/autostart/gablue-bigscreen-session.desktop
[Desktop Entry]
Type=Application
Name=Gablue Bigscreen Session Init
Comment=Initialise la session native Plasma Bigscreen (blacklist + mirroring écrans)
NoDisplay=true
Exec=/usr/libexec/gablue-bigscreen-session-init
Terminal=false
OnlyShowIn=KDE;
